# proxy/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs the NGINX daemon
# ==============================================================================

bashio::log.info "Starting NGINX..."
exec nginx
